# 93-keypad.rules
# Copy to /etc/udev/rules.d
ACTION=="add", SUBSYSTEM=="input", RUN+="/bin/systemctl restart keypad_client"
